# salt/base/time/ntp/files/etc/rc.local
#!/bin/sh -e

# Fix the proxmox problem without use capability
ntpd -g
